# packages/rootfs/reference/playwright-runner/rootfs/playwright-runner-ttl.timer
[Unit]
Description=Maximum lifetime for an ephemeral Playwright runner

[Timer]
OnActiveSec=45min
AccuracySec=1min
Unit=playwright-runner-ttl.service

[Install]
WantedBy=timers.target
