# packages/rootfs/reference/playwright-runner/rootfs/playwright-runner-ttl.service
[Unit]
Description=Stop an expired ephemeral Playwright runner

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
